# deploy/twd-order-backup.timer
[Unit]
Description=Daily TWD order database backup

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
